# puppet/manifests/andrewlorente.conf
start on startup
setuid andrewlorente
chdir /u/apps/andrewlorente/current
exec dist/build/andrewlorente/andrewlorente -p 8000
