# puppet/manifests/bloge.conf
start on startup
setuid bloge
chdir /u/apps/bloge/current
exec dist/build/bloge/bloge -p 8001 /u/apps/bloge/current/posts
