# roles/web-server/templates/bloge.conf.j2
description "bloge"

start on startup
setuid bloge
chdir /home/bloge/current
exec dist/build/bloge/bloge -p {{ item.0.port }} /home/bloge/current/posts
