# roles/web-server/templates/www.conf.j2
description "www"

start on startup
setuid www
chdir /home/www/current
exec dist/build/andrewlorente/andrewlorente -p {{ item.0.port }}
