# roles/web-server/files/andrewlorente.conf
description "andrewlorente"

start on startup
setuid andrewlorente
chdir /home/andrewlorente/current
exec dist/build/andrewlorente/andrewlorente -p 8000
