# roles/web-server/files/bloge.conf
description "bloge"

start on startup
setuid bloge
chdir /home/bloge/current
exec dist/build/bloge/bloge -p 8001 /home/bloge/current/posts
