# roles/web-server/files/identity.conf
description "identity"

start on (filesystem)
stop on runlevel [016]
respawn
console log
setuid identity
env HOME=/home/identity
env ENV=/home/identity/shared/.env
chdir /home/identity/current
exec ./Env/bin/gunicorn identity.app:app -b 0.0.0.0:5001 -w 3 --error-logfile /home/identity/shared/log/error.log --access-logfile /home/identity/shared/log/access.log
