# roles/web-server/files/paste.conf
description "paste"

start on (filesystem)
stop on runlevel [016]
respawn
console log
setuid paste
env HOME=/home/paste
chdir /u/apps/paste/current
exec node server.js
